# daemon/cmis.mount.in
[Mount]
Type=cmis
Exec=@libexecdir@/gvfsd-cmis
AutoMount=false
Scheme=cmis
